# scripts/systemd/rc-car.service
# Cloud RC car web/control/stream server.
# Install:
#   sudo cp scripts/systemd/rc-car.service /etc/systemd/system/
#   sudo cp server/.env.example /etc/rc-car.env   # then edit it
#   sudo systemctl enable --now rc-car
[Unit]
Description=Cloud RC Car server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/cloud-rc-car/server
EnvironmentFile=-/etc/rc-car.env
ExecStart=/usr/bin/python3 /home/pi/cloud-rc-car/server/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
